# data/user/cgroupify.service
[Unit]
Wants=cgroupify@%n.service
Before=cgroupify@%n.service

[Service]
Delegate=yes
