# data/user/cgroupify@.service.in
[Unit]
StopWhenUnneeded=yes
CollectMode=inactive-or-failed
ConditionControlGroupController=v2

[Service]
Type=simple
ExecStart=@libexecdir@/cgroupify %i
TimeoutStopSec=5
