# data/user/move-into-app-slice.service
[Service]
Slice=app.slice
# Enables cpu/io controllers
CPUWeight=100
IOWeight=100