# data/user/move-into-session-slice.service
[Service]
Slice=session.slice
# Disable memory controller
DisableControllers=memory
