# data/user/uresourced.service.in
[Unit]
Description=User resource assignment daemon
Requires=app.slice
After=pipewire.socket
Wants=pipewire.socket

After=app.slice

Before=graphical-session.target
PartOf=graphical-session.target
ConditionControlGroupController=v2

[Service]
# This daemon needs to be in the root slice so that the system daemon can
# detect it easily by checking whether the cgroup path exists.
Slice=-.slice

Type=notify
ExecStart=@libexecdir@/uresourced --user
TimeoutStopSec=5s


RestrictAddressFamilies=AF_UNIX
