# data/system/uresourced.service.in
[Unit]
Description=User resource assignment daemon
ConditionControlGroupController=v2

[Service]
Type=dbus
BusName=org.freedesktop.UResourced
ExecStart=@libexecdir@/uresourced

ProtectSystem=strict

# The daemon writes drop-in configurations for the user systemd instance
ReadWritePaths=/run/systemd
ReadOnlyPaths=/sys/fs/cgroup

# We only need DBus access
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=user@.service
