# data/system/user@.service
[Unit]
After=uresourced.service

[Service]
Delegate=cpu io memory
